# install/crontab-webpage-update-entry.cron

# cron entry to upload plot files to www
00 08 * * * cp -r /home/pi/UKRAA_PicoMuon/temp /var/www/html/
