# update/crontab-update.cron

# Start UKRAA PicoMuon crontab entry

# cron entry to process yesterdays raw muon counts per minute
30 00 * * * su pi -c "sudo bash $HOME/UKRAA_PicoMuon/scripts/processCPM.sh  >> $HOME/UKRAA_PicoMuon/logfiles/log-PicoMuonACM0.txt 2>>$HOME/UKRAA_PicoMuon/logfiles/log-error.txt"

# cron entry to process yesterdays raw muon adc values
15 03 * * * su pi -c "sudo bash $HOME/UKRAA_PicoMuon/scripts/processADC.sh  >> $HOME/UKRAA_PicoMuon/logfiles/log-PicoMuonACM0.txt 2>>$HOME/UKRAA_PicoMuon/logfiles/log-error.txt"

# cron entry to get NMDB data and process day, week and month data
30 04 * * * su pi -c "sudo bash $HOME/UKRAA_PicoMuon/scripts/processData.sh >> $HOME/UKRAA_PicoMuon/logfiles/log-PicoMuonACM0.txt 2>>$HOME/UKRAA_PicoMuon/logfiles/log-error.txt"

# cron entry to manage log files - on a Sunday 00:00 each week
00 00 * * 00 bash $HOME/UKRAA_PicoMuon/scripts/managelogfiles.sh

# cron entry to upload plot files to local website on RPi
30 09 * * * cp -r $HOME/UKRAA_PicoMuon/temp /var/www/html/

# End UKRAA PicoMuon crontab entry
