# update/crontab-update.cron
# cron entry to get neutron data from NMDB NEST
30 03 * * * su pi -c "/usr/bin/python3 /home/pi/UKRAA_PicoMuon/scripts/GetDataNeutron.py       >> /home/pi/UKRAA_PicoMuon/logfiles/log-PicoMuonACM0.txt 2>>/home/pi/UKRAA_PicoMuon/logfiles/log-error.txt"

# cron entry to process yesterdays raw muon counts per minute
30 00 * * * su pi -c "/usr/bin/python3 /home/pi/UKRAA_PicoMuon/scripts/ProcessDataCpmACM0.py   >> /home/pi/UKRAA_PicoMuon/logfiles/log-PicoMuonACM0.txt 2>>/home/pi/UKRAA_PicoMuon/logfiles/log-error.txt"

# cron entry to process yesterdays raw muon adc values
15 03 * * * su pi -c "/usr/bin/python3 /home/pi/UKRAA_PicoMuon/scripts/ProcessDataAdcACM0.py   >> /home/pi/UKRAA_PicoMuon/logfiles/log-PicoMuonACM0.txt 2>>/home/pi/UKRAA_PicoMuon/logfiles/log-error.txt"

# cron entry to process day, week and month data
30 04 * * * su pi -c "sudo bash /home/pi/UKRAA_PicoMuon/scripts/process.sh                     >> /home/pi/UKRAA_PicoMuon/logfiles/log-PicoMuonACM0.txt 2>>/home/pi/UKRAA_PicoMuon/logfiles/log-error.txt"

# cron entry to plot yesterdays various graphs
15 09 * * * su pi -c "sudo bash /home/pi/UKRAA_PicoMuon/scripts/plots.sh                       >> /home/pi/UKRAA_PicoMuon/logfiles/log-PicoMuonACM0.txt 2>>/home/pi/UKRAA_PicoMuon/logfiles/log-error.txt"

# cron entry to manage log files - on a Sunday 00:00 each week
00 00 * * 00 bash /home/pi/UKRAA_PicoMuon/scripts/managelogfiles.sh

# cron entry to upload plot files to local website on RPi
30 09 * * * cp -r /home/pi/UKRAA_PicoMuon/temp /var/www/html/